# docs/examples/hypervisors/kvm/kvm-baremetal/modules/kvm-installer-baremetal/scripts/secondary_vnics.service
[Unit]
Description=Enable Secondary vnics

[Service]
Type=oneshot
ExecStart=/usr/bin/configure-secondary-vnics.sh

[Install]
WantedBy=default.target
